# systemd/networkd-prefix-changed.target
[Unit]
Description=Network prefix changed trigger target
StopWhenUnneeded=yes
